# epaper-ai.service
[Unit]
Description=Générateur d'images IA pour e-Paper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/projects/mon-epaper/examples
EnvironmentFile=/home/pi/projects/mon-epaper/.env
ExecStart=/home/pi/.local/bin/uv run python3 /home/pi/projects/mon-epaper/examples/ai_image.py --mqtt
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
